# build/etc/supervisor.d/fixnagios.ini
[program:finishStartupScript]
command=/usr/local/bin/nagios.finishStartupScript
user=root
priority=5
autostart=false
startsecs=0
startretries=0

