# build/etc/supervisor.d/fnagios.ini
[program:nagios]
command=/usr/local/bin/fnagios
user=root
priority=5
autostart=false
startsecs=1

