# build/etc/supervisor.d/nagios.ini
[program:nagios]
command=/usr/local/nagios/bin/nagios /usr/local/nagios/etc/nagios.cfg
#command=/etc/init.d/nagios start
user=root
priority=5
autostart=true
autorestart=unexpected
startsecs=1
startretries=3
#exitcodes=0,2
stopsignal=QUIT
stopasgroup=false
killasgroup=false
stopwaitsecs=10
redirect_stderr=true
stdout_logfile=/var/log/nagios.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
stdout_events_enabled=false
serverurl=AUTO
