# etc/supervisor.d/nginx.ini
[program:nginx]
command=nginx -c /etc/nginx.conf
autostart=true
autorestart=true
priority=2
startretries=3
stopwaitsecs=10
stdout_logfile=/var/log/stdout_logfile.1.log
stdout_logfile_maxbytes=0
stderr_logfile=/var/log/stderr_logfile.1.log
stderr_logfile_maxbytes=0

